# etc/systemd/user/autostart-exam.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/lernstick-rdp-start-exam
