# etc/systemd/user/autostart-exam.service
[Service]
ExecStart=/usr/bin/autostart-exam
